# deploy/article-pipeline.example.service
# Article Pipeline — example systemd service
# Replace: /opt/article-pipeline → your actual install path
#          deploy-user → your actual user
#
# Setup:
#   sudo cp this file /etc/systemd/system/article-pipeline.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable article-pipeline
#   sudo systemctl start article-pipeline

[Unit]
Description=Article Pipeline UI
After=network.target

[Service]
Type=simple
User=deploy-user
WorkingDirectory=/opt/article-pipeline

# Loads credentials from .env in the project root
EnvironmentFile=/opt/article-pipeline/.env

ExecStart=/usr/bin/node ui/server.js
Restart=always
RestartSec=10

# Basic hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/article-pipeline

[Install]
WantedBy=multi-user.target
